# .config/autostart/redshift.desktop
[Desktop Entry]

Type=Application

Name=redshift

Exec=redshift
